# sync-docs.service
[Unit]
Description=Sincronizacion Docs-as-Code a Obsidian
After=default.target

[Service]
Type=simple
ExecStart=%h/scripts/sync_docs.sh
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Entorno limpio
Environment=LANG=C.UTF-8

[Install]
WantedBy=default.target
